# docker/run/fs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
exitcodes=0,2
directory=/

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:run_sshd]
command=/usr/sbin/sshd -D
environment=
stopwaitsecs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true

[program:run_cron]
command=/usr/sbin/cron -f
environment=
stopwaitsecs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true

[program:run_searxng]
command=/exe/run_searxng.sh
environment=SEARXNG_SETTINGS_PATH=/etc/searxng/settings.yml
user=searxng
directory=/usr/local/searxng/searxng-src
stopwaitsecs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true

[program:run_ui]
priority=10
command=/exe/run_A0.sh
environment=
user=root
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true

[program:run_tunnel_api]
command=/exe/run_tunnel_api.sh
environment=
user=root
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true

[eventlistener:the_listener]
command=python3 /exe/supervisor_event_listener.py
events=PROCESS_STATE_FATAL
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xvfb]
priority=5
command=/usr/bin/Xvfb :99 -screen 0 1280x768x24
environment=DISPLAY=:99
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:electron_shell]
priority=20
command=bash -c "while [ ! -f /a0/desktop/electron/run_electron.sh ]; do sleep 1; done; chmod +x /a0/desktop/electron/run_electron.sh; /a0/desktop/electron/run_electron.sh"
environment=DISPLAY=:99,NODE_ENV=production
user=root
stopwaitsecs=60
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
